# scripts/systemd/prompt-watchdog.timer
[Unit]
Description=Run the Marveen prompt watchdog every 30 seconds

[Timer]
OnBootSec=1min
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
